# dist/g11gkeys.service
[Unit]
Description=Handle Logitech G11 keys and call scripts for them

[Service]
Type=simple
ExecStart=/usr/bin/g11gkeys
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
